# etc/NetworkManager/dispatcher.d/01-flush-dhcp-client-state.sh
/etc/NetworkManager/dispatcher.d/pre-up.d/01-flush-dhcp-client-state.sh