# etc/NetworkManager/dispatcher.d/01-flush-dhcp-client-state.sh
/etc/NetworkManager/dispatcher.d/down.d/01-flush-dhcp-client-state.sh