# etc/NetworkManager/dispatcher.d/02-no-send-hostname.sh
/etc/NetworkManager/dispatcher.d/pre-up.d/02-no-send-hostname.sh