# etc/profile.d/hardened_malloc.sh
#!/usr/bin/sh
export LD_PRELOAD=/usr/lib/libhardened_malloc.so
